# systemd/opencode-bench.service
[Unit]
Description=opencode model/agent benchmark campaign (matrix x3 repeats)
# The campaign gates itself on the SAIA budget and can run for many hours;
# campaign.sh's flock makes accidental overlap harmless.

[Service]
Type=oneshot
WorkingDirectory=%h/benchmark-coding-agents
ExecStart=%h/benchmark-coding-agents/campaign.sh
# Campaigns legitimately run for hours; never let systemd kill a run mid-flight.
TimeoutStartSec=infinity
Nice=10
